# systemd/dwu.service
[Unit]
Description=DWU wallpaper updater

[Service]
Type=oneshot
ExecStart=/usr/bin/dwu --today
